# com.hoamdosa.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hoamdosa.telegram-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/JARVIS/Blanc Core/apps/hoamdosa/venv/bin/python3</string>
        <string>/Users/JARVIS/Blanc Core/apps/hoamdosa/telegram_bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/JARVIS/Blanc Core/apps/hoamdosa</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/JARVIS/Blanc Core/apps/hoamdosa/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/JARVIS/Blanc Core/apps/hoamdosa/bot_error.log</string>
</dict>
</plist>
